# systemd/nusantara-hw-detect.service
[Unit]
Description=NusantaraOS Hardware Detection (GPU Auto-Config)
Documentation=https://github.com/haikal-devtech/NusantaraOS
# Jalan sebelum display manager — PRD 5.2
# DefaultDependencies=no agar tidak inherit After=sysinit.target dll
# yang bisa bikin cycle dengan display-manager
DefaultDependencies=no
After=sysinit.target
Before=display-manager.service
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root
ExecStart=/bin/bash /usr/lib/nusantara/gpu-automation/hw-detect.sh
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nusantara-hw-detect

[Install]
WantedBy=sysinit.target
