# systemd/nusantara-hw-detect.service
[Unit]
Description=NusantaraOS Hardware Detection (GPU Auto-Config)
Documentation=https://github.com/haikal-devtech/nusantara-os
# Jalankan sebelum display manager — PRD 5.2 requirement
# "Run hw-detect service before display manager starts (systemd ordering: After=basic.target)"
After=basic.target
Before=display-manager.service
DefaultDependencies=no

[Service]
Type=oneshot
# Jalankan satu kali saat boot, tidak perlu restart
RemainAfterExit=yes
User=root
Group=root

ExecStart=/bin/bash /usr/lib/nusantara/gpu-automation/hw-detect.sh

# Timeout — PRD 5.2: fallback dalam < 30 detik
TimeoutStartSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nusantara-hw-detect

[Install]
WantedBy=multi-user.target
