# systemd/nusantara-guardian.service
[Unit]
Description=NusantaraOS System Guardian Daemon
Documentation=https://github.com/haikal-devtech/NusantaraOS
# Jalan setelah multi-user.target — tidak After=graphical.target
# karena itu yang bikin ordering cycle dengan display-manager/hw-detect
After=multi-user.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/lib/nusantara/guardian
ExecStart=/usr/bin/python3 /usr/lib/nusantara/guardian/main.py

# Restart otomatis kalau crash — PRD 5.3
Restart=on-failure
RestartSec=5

# Resource limit sesuai PRD 5.3: <50MB RAM, CPU idle
MemoryMax=50M
CPUQuota=10%

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nusantara-guardian

NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
